# deploy/webcodex-agent.service
[Unit]
Description=WebCodex local MCP agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=webcodex
Group=webcodex
Environment=HOME=/home/webcodex
Environment=CODEX_HOME=/home/webcodex/.codex
EnvironmentFile=/etc/webcodex/agent.env
ExecStart=/usr/local/bin/webcodex-agent
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
