# deploy/webcodex-gate.service
[Unit]
Description=WebCodex MCP gate
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=webcodex
Group=webcodex
Environment=HOME=/home/webcodex
EnvironmentFile=/etc/webcodex/gate.env
ExecStart=/usr/local/bin/webcodex-gate
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
